# net-proxy/make-proxy/files/mp-client
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

config="/etc/make-proxy/client"
pidfile="/run/make-proxy/mp-client.pid"
command="/usr/bin/erl -pa /opt/make-proxy/client/ebin -config ${config} +K true -s make_proxy_client start -noshell"

depend() {
    need net
}

start() {
    ebegin "Starting make-proxy-client"

    start-stop-daemon --start --exec ${command} \
    -- -c ${config} -f ${pidfile} >/dev/null 2>&1 &

    eend $?
}

stop() {
    ebegin "Stopping make-proxy-client"

    start-stop-daemon --stop --pidfile ${pidfile}

    eend $?
}
