# net-proxy/squid/files/squid.socket
[Unit]
Description=A full-featured web proxy cache

[Socket]

ListenStream=3128
Transparent=true

[Install]
WantedBy=sockets.target

