# app-misc/flox-bin/files/nix-daemon.initd
#!/sbin/openrc-run

description="Nix package manager daemon"
command="/usr/sbin/nix-daemon"
command_args="--daemon"
command_background="yes"
pidfile="/run/nix/nix-daemon.pid"
export NIX_SSL_CERT_FILE="@NIX_SSL_CERT_FILE@"

depend() {
	need localmount
	use net
}

start_pre() {
	checkpath -d -m 0755 /run/nix /nix/var/nix/daemon-socket
}
