# www-servers/bashttpd/files/http-alt.socket
[Unit]
Description=very simple HTTP index daemon


[Socket]
Accept=yes

ListenStream=8000
[Install]
WantedBy=sockets.target


