# www-servers/bashttpd/files/http@.service
[Service]

User=www
Group=root
WorkingDirectory=/usr/portage
StandardInput=socket
StandardOutput=socket
ExecStart=-/usr/bin/bashttpd
