# www-servers/bashttpd/files/http.socket
[Unit]
Description=very simple HTTP index daemon


[Socket]
Accept=yes

ListenStream=80
[Install]
WantedBy=sockets.target


