# net-proxy/shadowsocks-libev/files/shadowsocks-initd
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

ss_config="/etc/shadowsocks.json"

ss_pidfile="/run/shadowsocks.pid"
ss_command="/usr/bin/ss-local"

depend() {
    need net
}

checkconfig() {
    if [ ! -f ${ss_config} ]; then
        ewarn "${ss_config} does not exist."
    fi

    return 0
}

start() {
    checkconfig || return 1

    ebegin "Starting Shadowsocks"
    start-stop-daemon --start --exec ${ss_command} \
    -- -c ${ss_config} -f ${ss_pidfile} >/dev/null 2>&1 &
    eend $?
}

stop() {
    ebegin "Stopping Shadowsocks"
    start-stop-daemon --stop --pidfile ${ss_pidfile}
    eend $?
}

