# sys-apps/sysd/files/sysd.service
[Unit]
Description=the daemon who supplies firsthand system data

[Service]
ExecStart=/usr/sbin/sysd

[Install]
WantedBy=multi-user.target
