# sys-devel/distcc/files/3.0/init
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-devel/distcc/files/3.0/init,v 1.2 2008/10/25 06:39:17 matsuu Exp $

depend() {
	need net
	use avahi-daemon ypbind
}

start() {
	ebegin "Starting distccd"
	PATH="$(gcc-config --get-bin-path):${PATH}" \
	start-stop-daemon --start --quiet --exec "${DISTCCD_EXEC}" -- \
		--daemon --pid-file "${DISTCCD_PIDFILE}" --user distcc \
		${DISTCCD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "${DISTCCD_PIDFILE}"
	eend $?
}
