# net-proxy/goagent/files/goagent-initd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

GOAGENT_PID="--pidfile /var/run/goagent.pid"
GOAGENT_LOG="/var/log/goagent.log"

depend() {
    need net
    use logger
}

start() {
    ebegin "Starting GoAgent"
    start-stop-daemon --start --name goagent \
    $GOAGENT_PID --make-pid \
    --exec python2 /opt/goagent/local/proxy.py > $GOAGENT_LOG 2>&1 &
    eend $?
}

stop() {
    ebegin "Stopping GoAgent"
    start-stop-daemon --stop --name goagent $GOAGENT_PID \
    --exec python2 /opt/goagent/local/proxy.py
    eend $?
}

