# net-proxy/clash-verge-bin/files/clash-verge.initd
#!/sbin/openrc-run

name="clash-verge"
description="Clash Verge Service"
supervisor="supervise-daemon"
command="/opt/clash-verge/bin/clash-verge-service"
command_user="root:root"
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
    need net
    after bootmisc
}
