# net-proxy/clash/files/clash.initd
#!/sbin/openrc-run
# Copyright 2023 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

pidfile="/run/clash.pid"
command="/usr/bin/clash"
command_args="${CLASH_OPTS}"
command_background="yes"

depend() {
    need net
}
