# net-fs/ksmbd-tools/files/ksmbd.service
[Unit]
Description=ksmbd userspace daemon
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
User=root
Group=root
RemainAfterExit=yes
ExecStartPre=-/sbin/modprobe ksmbd
ExecStart=/usr/sbin/ksmbd.mountd -s
ExecReload=/usr/sbin/ksmbd.control -s && /usr/sbin/ksmbd.mountd
ExecStop=/usr/sbin/ksmbd.control -s

[Install]
WantedBy=multi-user.target
