# net-proxy/clash-for-windows-bin/files/clash-core-service.initd
#!/sbin/openrc-run
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Clash core service"
description="Clash core service created by Clash for Windows"
command="/opt/clash-for-windows-bin/resources/static/files/linux/x64/service/clash-core-service"
command_background="yes"
pidfile="/run/clash-core-service.pid"

depend() {
    need net
}
