# net-dns/phddns/files/init.d.phddns
#!/sbin/runscript

description="Phddns daemon"
depend() {
	need dbus
}

start() {
        phddns -d
}

stop() {
        killall phddns
}
