# sys-apps/aptly/files/aptly.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/aptly"
name="aptly"
description="aptly"
command_user="${APTLY_USER:-aptly}:${APTLY_GROUP:-aptly}"
command_args="serve ${OPTS}"

depend() {
	need net
}
